# deploy/roles/wadsworth/templates/wadsworth.service.j2
[Unit]
Description=Wadsworth Service

[Service]
ExecStart={{ start_script }}

[Install]
WantedBy=multi-user.target
