# deploy/roles/wadsworth/templates/wadsworth.service.j2
[Unit]
Description=Wadsworth Service

[Service]
ExecStart={{ wadsworth_install_path }}/run.sh

[Install]
WantedBy=multi-user.target
