# deploy/roles/wadsworth/templates/wadsworth.service.j2
[Unit]
Description=Wadsworth Service

[Service]
ExecStart=sh -c '{{ start_script }}'

[Install]
WantedBy=multi-user.target
